# slim/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0
childlogdir=/tmp

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ================================================================
# Startup order (via priority):
#   10  PostgreSQL
#   20  Temporal
#   30  control-api, data-api  (Python backends)
#   40  us, auth               (Node.js — preflight.sh waits for PG)
#   50  meta-manager           (depends on us + temporal)
#   60  ui-api                 (depends on us)
#   70  ui (node)              (depends on us, auth, meta-manager)
#   80  nginx                  (wraps ui node on 3030)
#
# Static per-service vars are in environment= below.
# Secrets (POSTGRES_PASSWORD_*, *_TOKEN, *_KEY, etc.) are exported
# by bootstrap.sh before exec supervisord and inherited by all programs.
# ================================================================

# ----------------------------------------------------------------
# PostgreSQL
# ----------------------------------------------------------------
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /data/pgdata -c logging_collector=off -c log_destination=stderr
user=postgres
priority=10
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# Temporal workflow engine
# ----------------------------------------------------------------
[program:temporal]
command=/etc/temporal/init-temporal.sh
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=5
environment=
    POSTGRES_USER="%(ENV_POSTGRES_USER_TEMPORAL)s",
    POSTGRES_PASSWORD="%(ENV_POSTGRES_PASSWORD_TEMPORAL)s",
    POSTGRES_DB="%(ENV_POSTGRES_DB_TEMPORAL)s",
    POSTGRES_DB_VISIBILITY="%(ENV_POSTGRES_DB_TEMPORAL_VISIBILITY)s",
    POSTGRES_TLS_ENABLED="false",
    BIND_ON_IP="127.0.0.1"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# control-api (Python/uWSGI, port 8001)
# ----------------------------------------------------------------
[program:control-api]
command=/opt/services/control-api/venv/bin/uwsgi
    --ini /opt/services/control-api/code/app/uwsgi/uwsgi-dl-api.ini
    --http-socket 0.0.0.0:8001
    --uid 0 --gid 0
directory=/opt/services/control-api
priority=30
autostart=true
autorestart=unexpected
startretries=10
startsecs=5
environment=
    VIRTUAL_ENV="/opt/services/control-api/venv",
    PATH="/opt/services/control-api/venv/bin:%(ENV_PATH)s",
    PYTHONHOME="",
    CONFIG_PATH="/opt/services/control-api/etc/backend-configs/api.yaml",
    HTTP_BIND_PORT="8001",
    DL_CRY_ACTUAL_KEY_ID="KEY",
    DL_CRY_KEY_VAL_ID_KEY="%(ENV_CONTROL_API_CRYPTO_KEY)s",
    US_HOST="http://localhost:8003",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    AUTH__TYPE="NATIVE",
    AUTH__JWT_ALGORITHM="PS256",
    AUTH__JWT_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s",
    RQE_FORCE_OFF="1",
    RQE_SECRET_KEY="",
    BI_API_UWSGI_WORKERS_COUNT="4",
    CONNECTOR_AVAILABILITY_VISIBLE="clickhouse,postgres,chyt,ydb,mysql,greenplum,mssql,oracle,trino,appmetrica_api,metrika_api"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# data-api (Python/Gunicorn+aiohttp, port 8002)
# ----------------------------------------------------------------
[program:data-api]
command=/opt/services/data-api/venv/bin/python -m gunicorn dl_data_api.app:create_gunicorn_app
    --bind 0.0.0.0:8002
    --workers 5
    --worker-class aiohttp.GunicornWebWorker
    --max-requests 5000
    --max-requests-jitter 500
    --timeout 90
    --keepalive 80
directory=/opt/services/data-api
priority=30
autostart=true
autorestart=unexpected
startretries=10
startsecs=5
environment=
    VIRTUAL_ENV="/opt/services/data-api/venv",
    PATH="/opt/services/data-api/venv/bin:%(ENV_PATH)s",
    PYTHONHOME="",
    CONFIG_PATH="/opt/services/data-api/etc/backend-configs/api.yaml",
    DL_CRY_ACTUAL_KEY_ID="KEY",
    DL_CRY_KEY_VAL_ID_KEY="%(ENV_CONTROL_API_CRYPTO_KEY)s",
    US_HOST="http://localhost:8003",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    AUTH__TYPE="NATIVE",
    AUTH__JWT_ALGORITHM="PS256",
    AUTH__JWT_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s",
    RQE_FORCE_OFF="1",
    RQE_SECRET_KEY="",
    CACHES_ON="0",
    MUTATIONS_CACHES_ON="0",
    BI_COMPENG_PG_ON="1",
    BI_COMPENG_PG_URL="postgres://%(ENV_POSTGRES_USER_COMPENG)s:%(ENV_POSTGRES_PASSWORD_COMPENG)s@localhost:%(ENV_POSTGRES_PORT)s/%(ENV_POSTGRES_DB_COMPENG)s",
    GUNICORN_WORKERS_COUNT="5"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# us (Node.js, port 8003)
# ----------------------------------------------------------------
[program:us]
command=/opt/services/us/scripts/preflight.sh
directory=/opt/services/us
priority=40
autostart=true
autorestart=unexpected
startretries=20
startsecs=10
environment=
    APP_PORT="8003",
    APP_ENV="prod",
    SKIP_INSTALL_DB_EXTENSIONS="1",
    USE_DEMO_DATA="0",
    AUTH_ENABLED="true",
    AUTH_ENDPOINT="http://localhost:8004",
    US_SURPRESS_DB_STATUS_LOGS="true",
    MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    AUTH_MASTER_TOKEN="%(ENV_AUTH_MASTER_TOKEN)s",
    AUTH_TOKEN_PUBLIC_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s",
    POSTGRES_DSN_LIST="postgres://%(ENV_POSTGRES_USER_US)s:%(ENV_POSTGRES_PASSWORD_US)s@localhost:%(ENV_POSTGRES_PORT)s/%(ENV_POSTGRES_DB_US)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# auth (Node.js, port 8004)
# ----------------------------------------------------------------
[program:auth]
command=/opt/services/auth/scripts/preflight.sh
directory=/opt/services/auth
priority=40
autostart=true
autorestart=unexpected
startretries=20
startsecs=10
environment=
    APP_PORT="8004",
    APP_ENV="prod",
    SKIP_INSTALL_DB_EXTENSIONS="1",
    DISABLE_WILDCARD_COOKIE="true",
    AUTH_SIGNUP_DISABLED="false",
    US_SURPRESS_DB_STATUS_LOGS="true",
    MASTER_TOKEN="%(ENV_AUTH_MASTER_TOKEN)s",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    TOKEN_PRIVATE_KEY="%(ENV_AUTH_TOKEN_PRIVATE_KEY)s",
    TOKEN_PUBLIC_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s",
    AUTH_ADMIN_PASSWORD="%(ENV_AUTH_ADMIN_PASSWORD)s",
    POSTGRES_DSN_LIST="postgres://%(ENV_POSTGRES_USER_AUTH)s:%(ENV_POSTGRES_PASSWORD_AUTH)s@localhost:%(ENV_POSTGRES_PORT)s/%(ENV_POSTGRES_DB_AUTH)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# meta-manager (Node.js, port 8005)
# ----------------------------------------------------------------
[program:meta-manager]
command=/opt/services/meta-manager/scripts/preflight.sh
directory=/opt/services/meta-manager
priority=50
autostart=true
autorestart=unexpected
startretries=20
startsecs=10
environment=
    APP_PORT="8005",
    APP_ENV="prod",
    TEMPORAL_AUTH_ENABLED="false",
    TEMPORAL_ENDPOINT="localhost:7233",
    AUTH_ENABLED="true",
    US_ENDPOINT="http://localhost:8003",
    UI_API_ENDPOINT="http://localhost:8006",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    EXPORT_DATA_VERIFICATION_KEY="%(ENV_EXPORT_DATA_VERIFICATION_KEY)s",
    AUTH_TOKEN_PUBLIC_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s",
    POSTGRES_DSN_LIST="postgres://%(ENV_POSTGRES_USER_META_MANAGER)s:%(ENV_POSTGRES_PASSWORD_META_MANAGER)s@localhost:%(ENV_POSTGRES_PORT)s/%(ENV_POSTGRES_DB_META_MANAGER)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# ui-api — headless API mode (Node.js, port 8006)
# ----------------------------------------------------------------
[program:ui-api]
command=node /opt/services/ui/dist/run
directory=/opt/services/ui
priority=60
autostart=true
autorestart=unexpected
startretries=20
startsecs=5
environment=
    APP_PORT="8006",
    APP_MODE="api",
    APP_ENV="production",
    US_ENDPOINT="http://localhost:8003",
    BI_API_ENDPOINT="http://localhost:8001",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# ui — full mode (Node.js, port 3030; nginx proxies 8080 → 3030)
# ----------------------------------------------------------------
[program:ui]
command=node /opt/services/ui/dist/run
directory=/opt/services/ui
priority=70
autostart=true
autorestart=unexpected
startretries=20
startsecs=5
environment=
    APP_PORT="3030",
    APP_MODE="full",
    APP_ENV="production",
    AUTH_POLICY="disabled",
    HC="0",
    YANDEX_MAP_ENABLED="0",
    EXPORT_WORKBOOK_ENABLED="true",
    AUTH_ENABLED="true",
    AUTH_SIGNUP_DISABLED="false",
    US_ENDPOINT="http://localhost:8003",
    AUTH_ENDPOINT="http://localhost:8004",
    BI_API_ENDPOINT="http://localhost:8001",
    BI_DATA_ENDPOINT="http://localhost:8002",
    META_MANAGER_ENDPOINT="http://localhost:8005",
    US_MASTER_TOKEN="%(ENV_US_MASTER_TOKEN)s",
    AUTH_MASTER_TOKEN="%(ENV_AUTH_MASTER_TOKEN)s",
    AUTH_TOKEN_PUBLIC_KEY="%(ENV_AUTH_TOKEN_PUBLIC_KEY)s"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ----------------------------------------------------------------
# nginx — reverse proxy / static files for UI (port 8080)
# ----------------------------------------------------------------
[program:nginx]
command=/etc/nginx/entrypoint.sh
priority=80
autostart=true
autorestart=unexpected
startretries=5
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
